# systemd/nexbreak-egress@.service
[Unit]
Description=NexBreak egress channel %i (packaging + delivery)
After=network-online.target
Wants=network-online.target
# Allow rapid recycle when SRT clients drop (ffmpeg listener is single-session).
StartLimitIntervalSec=120
StartLimitBurst=30

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/nexbreak/bin/nexbreak-egress --service-name %i --db /var/lib/nexbreak/nexbreak.sqlite
Restart=always
RestartSec=1
User=nexbreak
Group=nexbreak
# Don't wait on a stuck SRT listener — SIGTERM then SIGKILL the cgroup quickly.
TimeoutStopSec=8
KillMode=control-group
KillSignal=SIGTERM
FinalKillSignal=SIGKILL
SendSIGKILL=yes

# Same isolation model as nexbreak-proc@ — reads whatever local feed the
# router currently assigns it. Restarting an egress instance to change its
# destination does not touch the processing service producing its feed.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/nexbreak /var/log/nexbreak
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
