# systemd/nexbreak-controller.service
[Unit]
Description=NexBreak controller (config API + routing + audit)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/nexbreak/bin/nexbreak-controller --db /var/lib/nexbreak/nexbreak.sqlite --host 127.0.0.1 --port 8787
Restart=always
RestartSec=2
User=nexbreak
Group=nexbreak
# Same runtime dir as proc@ — required under ProtectSystem=strict so Verify
# can create /run/nexbreak/scte (ReadWritePaths alone is not enough if the
# path is not a systemd RuntimeDirectory for this unit).
RuntimeDirectory=nexbreak
RuntimeDirectoryMode=0770
Environment=NEXBREAK_RUN_DIR=/run/nexbreak
# Connect to /run/nexbreak/proc-*.sock (created by nexbreak-proc@ instances)
ReadWritePaths=/var/lib/nexbreak /var/log/nexbreak /run/nexbreak
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
