# systemd/nexbreak-verify.service
[Unit]
Description=NexBreak SCTE Verify API (return-feed monitor)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/nexbreak/bin/nexbreak-verify --db /var/lib/nexbreak/nexbreak.sqlite --host 127.0.0.1 --port 8788
Restart=always
RestartSec=2
User=nexbreak
Group=nexbreak
# Owns SCTE watch pid/state under /var/lib/nexbreak/scte (not /run — spawned
# helpers hit EROFS on RuntimeDirectory after ProtectSystem=strict).
RuntimeDirectory=nexbreak
RuntimeDirectoryMode=0770
Environment=NEXBREAK_RUN_DIR=/run/nexbreak
Environment=NEXBREAK_DATA=/var/lib/nexbreak
Environment=NEXBREAK_SCTE_DIR=/var/lib/nexbreak/scte
Environment=NEXBREAK_VERIFY_PORT=8788
Environment=NEXBREAK_LOG_DIR=/var/log/nexbreak
ReadWritePaths=/var/lib/nexbreak /var/log/nexbreak /run/nexbreak
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
# Keep spawned scte-watch in this cgroup so stop/restart cleans it up.
KillMode=control-group
TimeoutStopSec=8
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
